# usr/share/autostart/ksuperkey2.desktop
[Desktop Entry]
Exec=ksuperkey
X-DBUS-StartupType=none
Type=Service
X-KDE-StartupNotify=false
OnlyShowIn=KDE;
X-KDE-autostart-phase=1
Hidden=true